# scripts/com.election-intelligence.platform.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.election-intelligence.platform</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_ROOT/scripts/start-platform.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>REPO_ROOT/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>REPO_ROOT/logs/launchd.err.log</string>
  <key>WorkingDirectory</key>
  <string>REPO_ROOT</string>
</dict>
</plist>
